# supervisord.conf
[include]
files = /home/atlas/supervisor.d/*.conf

[unix_http_server]
file=/var/run/supervisor.sock  ; No auth needed — socket is container-internal only

[supervisord]
nodaemon=true
logfile=/atlas/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/atlas/logs

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:init]
command=/atlas/app/init.sh
autorestart=false
startsecs=0
startretries=0
priority=10
stdout_logfile=/atlas/logs/init.log
stderr_logfile=/atlas/logs/init-error.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/atlas/logs/nginx.log
stderr_logfile=/atlas/logs/nginx-error.log

[program:inbox-mcp]
command=bun run /atlas/app/inbox-mcp/index.ts
autostart=false
autorestart=true
priority=30
stdout_logfile=/atlas/logs/inbox-mcp.log
stderr_logfile=/atlas/logs/inbox-mcp-error.log

[program:web-ui]
command=bun run /atlas/app/web-ui/index.ts
autostart=false
autorestart=true
priority=40
stdout_logfile=/atlas/logs/web-ui.log
stderr_logfile=/atlas/logs/web-ui-error.log

[program:watcher]
command=/atlas/app/watcher.sh
autostart=false
autorestart=true
priority=50
stdout_logfile=/atlas/logs/watcher.log
stderr_logfile=/atlas/logs/watcher-error.log

[program:supercronic]
command=supercronic /home/atlas/crontab
autostart=false
autorestart=true
priority=60
stdout_logfile=/atlas/logs/cron.log
stderr_logfile=/atlas/logs/cron-error.log

[program:playwright-mcp]
command=npx @playwright/mcp --headless --no-sandbox --port 8931 --host 127.0.0.1
autostart=false
autorestart=true
priority=36
stdout_logfile=/atlas/logs/playwright-mcp.log
stderr_logfile=/atlas/logs/playwright-mcp-error.log
environment=DISPLAY=""
