# supervisord.conf
[include]
files = /home/agent/supervisor.d/*.conf

[unix_http_server]
file=/var/run/supervisor.sock  ; No auth needed — socket is container-internal only

[supervisord]
nodaemon=true
logfile=/talos/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/talos/logs

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:init]
command=/talos/app/init.sh
autorestart=false
startsecs=0
startretries=0
priority=10
stdout_logfile=/talos/logs/init.log
stderr_logfile=/talos/logs/init-error.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/talos/logs/nginx.log
stderr_logfile=/talos/logs/nginx-error.log

[program:talos-mcp]
command=bun run /talos/app/talos-mcp/index.ts
autostart=false
autorestart=true
priority=30
stdout_logfile=/talos/logs/talos-mcp.log
stderr_logfile=/talos/logs/talos-mcp-error.log

[program:web-ui]
command=bun run /talos/app/web-ui/index.ts
autostart=false
autorestart=true
priority=40
stdout_logfile=/talos/logs/web-ui.log
stderr_logfile=/talos/logs/web-ui-error.log

[program:supercronic]
command=supercronic /home/agent/crontab
autostart=false
autorestart=true
priority=60
stdout_logfile=/talos/logs/cron.log
stderr_logfile=/talos/logs/cron-error.log
