# app/defaults/supervisor.d/cloudflared.conf
[program:cloudflared]
command=/home/agent/bin/cloudflared tunnel --url http://localhost:8080 --no-autoupdate
autostart=true
autorestart=true
stdout_logfile=/talos/logs/cloudflared.log
stderr_logfile=/talos/logs/cloudflared-error.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
