# app/defaults/supervisor.d/email-poller.conf
[program:email-poller]
command=/talos/app/scripts/email-poller.sh
autostart=false
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/home/agent/.logs/email-poller.log
stderr_logfile=/home/agent/.logs/email-poller.err
environment=OUTLOOK_TENANT_ID="%(ENV_OUTLOOK_TENANT_ID)s",OUTLOOK_CLIENT_ID="%(ENV_OUTLOOK_CLIENT_ID)s",OUTLOOK_CLIENT_SECRET="%(ENV_OUTLOOK_CLIENT_SECRET)s",OUTLOOK_MAILBOX="%(ENV_OUTLOOK_MAILBOX)s",OUTLOOK_ALLOWED_SENDERS="%(ENV_OUTLOOK_ALLOWED_SENDERS)s"
