# app/defaults/supervisor.d/stt-server.conf
[program:stt-server]
command=parakeet --port 5092 --host 127.0.0.1
autostart=false
autorestart=true
stdout_logfile=/atlas/logs/stt-server.log
stderr_logfile=/atlas/logs/stt-server-error.log
