# app/defaults/supervisor.d/tunnel-update.conf
[program:tunnel-update]
command=/talos/app/scripts/update-tunnel-url.sh
autostart=true
autorestart=false
startsecs=0
startretries=1
priority=90
stdout_logfile=/talos/logs/tunnel-update.log
stderr_logfile=/talos/logs/tunnel-update-error.log
