# app/defaults/crontab
# Agent Crontab (supercronic)
# Timezone via TZ env variable in the container

# Check for due reminders every minute
* * * * *  bun /talos/app/triggers/manage-reminders.ts check >> /talos/logs/reminders.log 2>&1

# Daily session cleanup at 06:00 (special: uses --max-turns and session resume)
0 6 * * *  /talos/app/triggers/cron/daily-cleanup.sh

# Watchdog: detect and kill stuck trigger-runner processes (>30min)
*/5 * * * *  /talos/app/scripts/watchdog-triggers.sh 2>/dev/null

# --- A2M-specific ---

# Anton calendar reminders (every minute)
* * * * *  /talos/app/scripts/anton-reminders.py >> /talos/logs/anton-reminders.log 2>&1

# Daemon health checks (every 5 minutes)
*/5 * * * *  /talos/app/scripts/check-daemon-health.sh >> /talos/logs/daemon-health.log 2>&1

